# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off; error_log /dev/stdout info;"
redirect_stderr=true
